# curator4-1618243501.conf
[program:curator4-1618243501]
command=python3 /home/curator4/curator4-1618243501/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1618243501/test.err.log
stdout_logfile=/home/curator4/curator4-1618243501/test.out.log